# docker/supervisord.conf
; ─────────────────────────────────────────────────────
;  Supervisord — Process Manager
;  Manages Nginx, Next.js, and Streamlit together
; ─────────────────────────────────────────────────────

[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid
loglevel=info

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=10

[program:nextjs]
command=/app/presentation/node_modules/.bin/next start -p 3000
directory=/app/presentation
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=NODE_ENV="production",PORT="3000"
priority=20

[program:streamlit]
command=streamlit run /app/streamlit_app.py
    --server.port=8501
    --server.headless=true
    --server.address=0.0.0.0
    --server.enableCORS=false
    --server.enableXsrfProtection=false
    --server.baseUrlPath=/demo
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=
    MLFLOW_TRACKING_URI="%(ENV_MLFLOW_TRACKING_URI)s",
    MLFLOW_S3_ENDPOINT_URL="%(ENV_MLFLOW_S3_ENDPOINT_URL)s",
    AWS_ACCESS_KEY_ID="%(ENV_AWS_ACCESS_KEY_ID)s",
    AWS_SECRET_ACCESS_KEY="%(ENV_AWS_SECRET_ACCESS_KEY)s",
    CLOUD_INFERENCE_URL="%(ENV_CLOUD_INFERENCE_URL)s"
priority=30
